# scripts/pinterest/com.paintcolorhq.pinterest-drip.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>__LABEL__</string>
  <!-- Invoke /bin/bash DIRECTLY (not the script via its shebang): macOS TCC
       attributes file access to the launchd program, so the FDA-granted
       /bin/bash must be the program — otherwise runs fail with "Operation not
       permitted" reading the repo (~/Documents) / images (~/Desktop). -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__DRIPSH__</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__NODE_DIR__:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <!-- Spread ~6 pins/day across 3 runs (9am / 1pm / 5pm EDT) instead of one
       burst — looks organic on Pinterest and catches more of the feed cycle. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardErrorPath</key><string>__REPO__/scripts/pinterest/drip.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
